# new_op_supervisor.ini
[program:operation_agent]
command=uvicorn main:app --host 0.0.0.0 --port 11011 --reload
directory=/root/xwh/new_operation/new_op_platform
autostart=true
startsecs=10
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile_maxbytes=20MB
stdout_logfile_backups = 1
stdout_logfile=/root/xwh/operation_service/out/operation_agent.stdout
stopasgroup=true
killasgroup=true
